# roles/pystore/files/pystore.socket
[Unit]
Description=pystore socket

[Socket]
ListenStream=/run/pystore.sock

[Install]
WantedBy=sockets.target

