# roles/gunicorn/files/gunicorn_ecommerce.socket
[Unit]
Description=gunicorn ecommerce socket

[Socket]
ListenStream=/run/gunicorn_ecommerce.sock

[Install]
WantedBy=sockets.target
